# kegbot.service
[Unit]
Description=Kegbot Base
After=network.target

[Service]
ExecStart=/usr/local/bin/kegbot
Type=simple

[Install]
WantedBy=default.target
